# systemd/solar_heating_web_gui.service
[Unit]
Description=Solar Heating System v3 - Web GUI
After=network.target solar_heating_api.service
Wants=solar_heating_api.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/sun_heat_and_ftx/python/v3
ExecStart=/usr/bin/python3 /home/pi/sun_heat_and_ftx/python/v3/web_server.py --host 0.0.0.0 --port 8080 --api-port 5001
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=solar_heating_web_gui

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp/solar_heating

[Install]
WantedBy=multi-user.target
